# libsnow-helper/data/libsnow-helper.service.in
[Unit]
Description=LibSnow Helper Daemon

[Service]
Type=dbus
BusName=org.snowflakeos.LibSnow.Helper1
ExecStart=@bindir@/libsnow-helper dbus
